# luci/luci-app-mosdns/root/etc/uci-defaults/luci-mosdns
#!/bin/sh
#
# luci-app-mosdns first-install wiring.
#   - always owns the configfile pointer (this app manages config_custom.yaml)
#   - enables the service + dnsmasq redirect ONCE (sentinel-guarded so later
#     LuCI toggles survive luci-app upgrades)
#   - registers ucitrack so UCI mosdns edits reload the init

SENTINEL=/etc/.luci_app_mosdns

uci set mosdns.config.configfile='/etc/mosdns/config_custom.yaml'

if [ ! -f "$SENTINEL" ]; then
	uci set mosdns.config.enabled='1'
	uci set mosdns.config.redirect='1'
	touch "$SENTINEL"
fi

uci -q commit mosdns

if [ -f /etc/config/ucitrack ]; then
	uci -q batch <<-EOF >/dev/null
		delete ucitrack.@mosdns[-1]
		add ucitrack mosdns
		set ucitrack.@mosdns[-1].init=mosdns
		commit ucitrack
	EOF
fi

# Boot-enroll the init unconditionally (writes the S99 rc.d symlink — safe even
# inside the image-build chroot). Start mosdns right now on a live device when
# it isn't already running, closing the "installed but not running until a
# reboot" gap without restarting a healthy service on upgrades.
/etc/init.d/mosdns enable >/dev/null 2>&1
if [ -z "$IPKG_INSTROOT" ] && ! pidof mosdns >/dev/null 2>&1; then
	/etc/init.d/mosdns start >/dev/null 2>&1
fi

rm -rf /tmp/luci-*
exit 0
